# system_files/usr/lib/systemd/system/nordvpnd.socket.d/10-nordvpn-group.conf
[Unit]
After=nordvpn-group.service
Requires=nordvpn-group.service
